# etc/init/udev-static.conf
# udev-static - creates static devices on /dev devtmpfs
#

description    "create static devices on /dev devtmpfs"

start on (started udev)

task

exec /bin/cp -ap /lib/udev/devices/* /dev/
